# backup-service/nexus-backup.service
[Unit]
Description=Nexus data backup

[Service]
Type=oneshot
User=nexus
ExecStart=/usr/local/bin/backup.sh